# UI/elastic.service
[Unit]
Description=Service that keeps running the echo-server from startup.